# com.replyscreener.afternoon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.replyscreener.afternoon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/manojjoshi/Desktop/scraper/run_random_window.sh</string>
        <string>afternoon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/manojjoshi/Desktop/scraper</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>15</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/manojjoshi/Desktop/scraper/logs/launchd-afternoon.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/manojjoshi/Desktop/scraper/logs/launchd-afternoon.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
